# provisioning/monitoring/envoy.service
[Unit]
Description=Envoy
Wants=network-online.target
After=network-online.target

[Service]
User=root
Restart=on-failure
ExecStart=/bin/envoy \
  --config-path /vagrant/provisioning/monitoring/envoy-me.yaml

[Install]
WantedBy=default.target
